# FoldingAtHome/supervisord.conf
[supervisord]
nodaemon=true
[program:fahclient]
command=/usr/bin/FAHClient --config /config/config.xml
